# files/fbi-slideshow.service
[Unit]
Description=Dynamic Long-Running Slideshow Service
# Test comment to trigger ansible change
# Ensures network and basic graphics are ready
After=network-online.target graphical.target
Wants=network-online.target

[Service]
User=root
Group=root
ExecStart=/home/tcarter/start-slideshow.sh

# Restart the service if it ever exits (crashes, is killed, etc.) 
Restart=always

# Wait 30 seconds before attempting a restart to prevent rapid-fire loops
RestartSec=30
WatchdogSec=90

# Note: These lines prevent script output from filling the system log.
# For debugging, you can comment them out to see logs via 'journalctl'.
StandardOutput=null
StandardError=null

[Install]
WantedBy=multi-user.target